# testdata/jobs1.cron
* * * * * foo
*/10 * * * * bar
# Just a comment
15 23 * * * /usr/local/bin/backup
