# testdata/jobs2.cron
* * * * * /bin/sleep 10
# Blank line follows:

* * * * * /bin/bogus
# Line with only whitespace follows:
    	 
* * * * * /bin/lastjob
